# deploy/scripts/crontab
# =============================================================================
# Crontab for PostgreSQL Backup Container
# =============================================================================
# GDPR Compliance: Automated daily backups at 02:00 UTC
# =============================================================================

# Daily backup at 02:00 UTC
0 2 * * * /usr/local/bin/backup.sh >> /var/log/backup.log 2>&1

# Note: Lifecycle rules in S3 handle retention:
# - Daily backups: 30 days
# - Weekly backups (Sunday): 90 days
# - Monthly backups (1st): 365 days
